# deploy/bluebex-web.service
[Unit]
Description=Bluebex Teams Web (Next.js)
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/opt/bluebex-teams/apps/web/.next/standalone
EnvironmentFile=/opt/bluebex-teams/.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
